# .config/systemd/user/display.target
[Unit]
Description=Graphical target
AllowIsolate=yes
