# .config/systemd/user/display.target
[Unit]
Description=Foobar boot target
AllowIsolate=yes
